# supervisord.conf
[supervisord]
nodaemon=true

[program:publisher]
command=bash -c "python /app/publisher.py"
autostart=true
autorestart=true
stdout_logfile=/var/log/publisher.log
stderr_logfile=/var/log/publisher.err.log



[program:fastapi]
command=bash -c "python /app/main.py --host 0.0.0.0 --port 8000"
autostart=true
autorestart=true
stdout_logfile=/var/log/fastapi.log
stderr_logfile=/var/log/fastapi.err.log
